# resources/work/crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.Jy1TjM9ksE installed on Wed Feb 10 09:59:47 2021)
# (Cron version -- $FreeBSD: src/usr.sbin/cron/crontab/crontab.c,v 1.24 2006/09/03 17:52:19 ru Exp $)
0 * * * * sh /Users/tevanc/Google\ Drive\ File\ Stream/My\ Drive/devStuff/scripts/random/oathCaller/wrapper.sh
#*/5 * * * * export AWS_DEFAULT_PROFILE=zgm-dev-developer && cd /Volumes/zDev/misc/gitlab_manipulation/dependency_mapper && source ./../venv/bin/activate && python3 dependency_visualizer.py
*/10 * * * * sh /Users/tevanc/beanieScraper.sh
